# source/etc/systemd/system/getty@.service.d/override.conf
[Service]
# Clear the default ExecStart line
ExecStart=
# Replace it with agetty in autologin mode for the user 'user'
ExecStart=-/sbin/agetty --autologin user --noclear %I $TERM
